# templates/upstart-runsvdir-start.conf.j2
start on runlevel [12345]
respawn
stop on shutdown
exec {{ runit_binary_dir}}/runsvdir-start
